# udev_rules/13-fix-keyboard.rules
SUBSYSTEMS=="input", AATRS{name}=="Keychron K3", RUN+="/home/matthew/scripts/fix_keyboard"
